# pxe-boot/do-docker-host/files/etc/knock/systemd/system/firewall.service
[Unit]
Description=Network Firewall Configuration Service
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/firewall start
ExecStop=/usr/local/sbin/firewall stop

[Install]
WantedBy=multi-user.target